# test-monitor-run.service
[Unit]
Description=Test Monitor Runner
After=network.target

[Service]
Type=oneshot
User=root
ExecStart=/usr/local/bin/test_monitor.sh
StandardOutput=journal